# etc/cron.d/duplicacy
# m h  dom mon dow   command
13 * * * * /usr/local/bin/duplicacy-wrapper.sh

